# deploy/systemd/irrigation.service.template
[Unit]
Description=Automated irrigation system
After=multi-user.target

[Service]
Type=simple
User=__RUN_USER__
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=-/etc/default/irrigation-system
ExecStart="__BINARY__" run
Restart=always
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
